# ssh-tunnel/crontab
SHELL=/bin/bash
* * * * * /opt/ssh_tunnel/connect_to_hosts.sh